# deploy/config_proposal_worker.service
[Unit]
Description=Humble Capital Config Proposal Worker
After=network.target

[Service]
Type=simple
User=root
Group=root
WorkingDirectory=/root/tradingbot
EnvironmentFile=/root/tradingbot/.env
ExecStart=/usr/bin/python3 /root/tradingbot/config_proposal_worker.py
Restart=always
RestartSec=10

[Install]
WantedBy=multi-user.target
